# deploy/sroolify-agent.service
[Unit]
Description=Sroolify Monitoring Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sroolify-agent --server ws://SERVER_HOST:8080/ws --token AGENT_TOKEN --agent-id AGENT_ID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
